# lin.service
[Unit]
Description=lin

[Service]
ExecStart=/home/lin/start.sh

[Install]
WantedBy=multi-user.target
