# linano.service
[Unit]
Description=linano

[Service]
ExecStart=/home/linano/start

[Install]
WantedBy=multi-user.target
